# infra/systemd/orders-sync.timer
[Unit]
Description=Run orders-sync every minute

[Timer]
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
EOF
