# infra/systemd/orders-sync.timer
[Unit]
Description=Run orders-sync every minute

[Timer]
# 부팅 30초 후 시작
OnBootSec=30s
# 직전 실행 시점 기준 60초 주기
OnUnitActiveSec=60s
# 타이머 지터: 여러 인스턴스일 때 폭주 방지(선택)
RandomizedDelaySec=5s
# 서버가 꺼져 있던 동안의 실행 시점을 보존(부팅시 즉시 한 번 실행)
Persistent=true
# 타이머가 호출할 서비스
Unit=orders-sync.service
# 타이머 정밀도(선택)
AccuracySec=5s

[Install]
WantedBy=timers.target
