# infra/systemd/orders-sync.service
[Unit]
Description=Sync Imweb Orders via internal endpoint
Wants=network-online.target
After=network-online.target premise.service

[Service]
Type=oneshot
User=prewoos2018
ExecStart=/usr/bin/curl --silent --show-error --fail-with-body --max-time 20 "http://127.0.0.1:5000/internal/orders-sync?sync_token=d2c0c892-68c3-4726-8553-ebba2c2ae928"
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
